# deploy/systemd/server-monitoring-celery.service
# NOTE: Redis service name varies by distro:
#   Debian/Ubuntu: redis-server.service
#   RHEL/Fedora/Arch: redis.service
[Unit]
Description=Server Monitoring Celery Worker
After=network.target redis-server.service
Requires=redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/server-monitoring
EnvironmentFile=/etc/server-monitoring/env
ExecStart=/opt/server-monitoring/.venv/bin/celery -A config worker \
    --loglevel=info \
    --concurrency=2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target